# crontab.cron
SHELL=/bin/bash
0 6 * * * cd .. && usr/local/bin/python /home/web-scraper/web_scraper.py
0 7 * * * cd .. && usr/local/bin/python /home/web-scraper/read_json.py
