# crontab
0 6 * * * python /web-scraper/web_scraper.py
