# crontab.cron
SHELL=/bin/bash
40 16 * * * cd .. && usr/local/bin/python /home/web-scraper/web_scraper.py
* * * * * echo "test" >> /tst.txt
